# deploy/meshcom-webdesk.service
[Unit]
Description=MeshCom WebDesk
After=network.target

[Service]
Type=simple
User=meshcom
WorkingDirectory=/opt/meshcom
ExecStart=/opt/meshcom/MeshcomWebDesk
Restart=on-failure
RestartSec=5

# Override settings via environment variables
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://0.0.0.0:5162
Environment=Meshcom__LogPath=/var/log/meshcom

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meshcom-webdesk

[Install]
WantedBy=multi-user.target
